# Linux/.config/systemd/user/waybar.service
[Unit]
Description=Waybar
PartOf=graphical-session.target
After=graphical-session.target

[Service]
ExecStart=waybar
Restart=on-failure

[Install]
WantedBy=graphical-session.target
